# systemd/proxy-sync.service
[Unit]
Description=Update proxy configuration based on GNOME settings
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/proxy-sync.sh monitor
Restart=on-failure
RestartSec=10

# Run as root to modify system config files
User=root
Group=root

# Environment
Environment=XDG_RUNTIME_DIR=/run/user/1000

# Security settings
ProtectSystem=false
PrivateDevices=false
ProtectHome=false
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
